# tf-k8ssrv/modules/vm/scripts/systemd/ansible-pull.service.tpl
[Unit]
Description=Run ansible-pull

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ansible-pull-runner.sh

